# meta-guacamayo/recipes-core/guacamayo-session/guacamayo-session-egl/raspberrypi/guacamayo-session-egl
#!/bin/sh
#
### BEGIN INIT INFO
# Provides: rygel
### END INIT INFO

USER=rygel
XDGAUTOSTART=/etc/xdg/autostart

. /etc/guacamayo-session-common

case "$1" in
  start)
       . /etc/profile
       ensurehome
       cleanupdbus
       setupdbus

       chown $USER /dev/tty[01234567]

       # We should sort this out in bbappend to vc-graphics, but the package
       # is still in flux
       chgrp video /dev/vchiq
       chmod g+rw /dev/vchiq

       export CLUTTER_INPUT_BACKEND=evdev
       # the g_slice allocator crashes on the RPi, malloc does not
       # though the device is really memory starved
       export G_SLICE=always-malloc
       # dump the environment for debugging purposes
       export > $HOME/.environment

       # if the gstreamer binary registry does not exist yet, generate it now
       # otherwise this will happen on demand when Mex first start and the
       # clutter plugin will fail and get blacklisted due to EGL context error
       GSTREG=`ls -1 $HOME/.gstreamer*/registry*.bin 2>/dev/null`
       if [ x$GSTREG = x ]; then
          echo "Introspecting Gstreamer"
          sudo -E -u $USER gst-inspect &> /dev/null
       fi

       echo "Starting Rygel"
       /sbin/start-stop-daemon -b -c $USER -S -x /usr/bin/rygel

       echo "Running XDG Autostarts"
       if [ -d $XDGAUTOSTART ]; then
	   for SCRIPT in $XDGAUTOSTART/*; do
               CMD=`grep Exec= $SCRIPT | sed -e s/.*Exec=// -e s/%.//g`
               sudo -b -E -u $USER -g $USER $CMD
	   done
       fi

       echo "Starting MEX"
       sudo -b -E -u $USER -g $USER /usr/bin/guacamayo-watchdog --debug=all -- media-explorer
  ;;

  stop)
	killmex
        killrygel
	killpulse
	killdbus
  ;;

  restart)
	$0 stop
        sleep 1
        $0 start
  ;;

  *)
        echo "usage: $0 { start | stop | restart }"
  ;;
esac

exit 0
