# meta-guacamayo/recipes-core/guacamayo-session/guacamayo-session-mediaserver/guacamayo-session-mediaserver
#!/bin/sh
#
### BEGIN INIT INFO
# Provides: rygel
### END INIT INFO

USER=rygel

. /etc/guacamayo-session-common

case "$1" in
  start)
       . /etc/profile
       ensurehome
       setupdbus

       echo "Starting Rygel"
       /sbin/start-stop-daemon -b -c $USER -S -x /usr/bin/rygel -- \
                               --plugin-option=MediaExport:uris:/usr/share/demos/,$HOME \
                               --title=MediaExport:"Rygel Media Server"
  ;;

  stop)
        killrygel
	killdbus
  ;;

  restart)
	$0 stop
        sleep 1
        $0 start
  ;;

  *)
        echo "usage: $0 { start | stop | restart }"
  ;;
esac

exit 0
