# meta-guacamayo/recipes-core/guacamayo-session/guacamayo-session-audioplayer/guacamayo-session-audioplayer
#!/bin/sh
#
### BEGIN INIT INFO
# Provides: rygel
### END INIT INFO

USER=rygel

. /etc/guacamayo-session-common

case "$1" in
  start)
       . /etc/profile
       ensurehome
       setupdbus

       # dump the environment for debugging purposes
       export > $HOME/.environment

       echo "Starting Pulse Audio server"
       # can't use start-stop-daemon for pa, it starts fine but does not work
       sudo -b -E -u $USER -g $USER /usr/bin/pulseaudio -D --exit-idle-time=-1

       echo "Starting Rygel"
       /sbin/start-stop-daemon -b -c $USER -S -x /usr/bin/rygel

       /usr/bin/guacamayo-watchdog --keep-alive --debug=all -- guacamayo-cli
  ;;

  stop)
        killrygel
	killpulse
	killdbus
  ;;

  restart)
	$0 stop
        sleep 1
        $0 start
  ;;

  *)
        echo "usage: $0 { start | stop | restart }"
  ;;
esac

exit 0
