# meta-guacamayo/recipes-support/guacamayo-watchdog/files/etc/xdg/autostart/guacamayo-watchdog-mex.desktop
[Desktop Entry]
Type=Application
Name=Guacamayo MEX watch dog
Exec=guacamayo-watchdog --debug=all -- media-explorer -f
